# etc/hotplug.d/iface/70-fap
#!/bin/sh
[ "$ACTION" = ifup ] && setupfap "$INTERFACE" | sh
