# etc/init.d/qos
#!/bin/sh /etc/rc.common
# Copyright (C) 2006 OpenWrt.org

START=50

boot() {
	/usr/lib/qos/generate.sh firewall | sh
	/usr/lib/qos/wireless.sh boot
}

start() {
	/usr/lib/qos/generate.sh firewall start | sh
	/usr/lib/qos/wireless.sh reload
}

stop() {
	/usr/lib/qos/generate.sh firewall stop | sh
	/usr/lib/qos/wireless.sh stop
}
