# usr/lib/opkg/info/cwmpd.list
/usr/lib/cwmpd/transfers/checkVoiceCalls.lua
/usr/lib/cwmpd/libuciplugin.so
/usr/lib/cwmpd/transfers/upgrade.sh
/usr/bin/cwmpd
/usr/lib/cwmpd/transfers/exec_config.lua
/etc/init.d/cwmpdboot
/usr/bin/cwmpevents
/etc/uci-defaults/70-cwmpd
/etc/init.d/cwmpd
/usr/lib/cwmpd/transfers/config.sh
/lib/functions/firewall-cwmpd.sh
/usr/lib/libcwmpTransformerBackend.so
/usr/lib/cwmpd/libtransferplugin.so
