# usr/lib/opkg/info/ethswctl.list
/etc/init.d/switchqospause
/usr/bin/ethswctl
