# usr/lib/opkg/info/lua-bcm.list
/usr/bin/setupfap
/etc/init.d/portmirror
/usr/bin/setportmirror
/etc/config/portmirror
/usr/bin/setxdsloemparam
/etc/hotplug.d/iface/70-fap
/usr/lib/lua/luabcm.so
