# usr/lib/opkg/info/mwan.list
/etc/init.d/mwan
/etc/config/mwan
/etc/hotplug.d/iface/15-mwan
/etc/ld.so.preload
/usr/lib/libmwan-preload.so
