# usr/lib/opkg/info/qos-tch.list
/etc/init.d/qos
/etc/init.d/dscp2pbit
/usr/lib/qos/wireless.sh
/etc/hotplug.d/iface/05-interface-remarking
/etc/hotplug.d/iface/00-qos
/usr/lib/qos/generate.sh
/etc/hotplug.d/iface/10-dscp2pbit
/usr/lib/qos/interface-remarking.sh
/etc/config/qos
/etc/config/dscp2pbit
